# autofill-watcher.service
[Unit]
Description=Sunbiz Autofill Watcher
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/company-questionnaire
Environment="DISPLAY=:1"
EnvironmentFile=/home/ubuntu/.airtable_env
ExecStart=/usr/bin/python3 /home/ubuntu/company-questionnaire/autofill_watcher.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

